# software/systemd/acs-endpoint-provisioning.service
[Unit]
Description=Provisioning Modbus devices configuration on startup

[Service]
Type=simple
ExecStart=/usr/bin/python /home/acs/app/acs/software/provision/provision.py
WorkingDirectory=/home/acs/app/acs/software/provision/
User=acs

[Install]
WantedBy=multi-user.target

